# deploy/hermes-agent-web-chat-backend.service
[Unit]
Description=Hermes Agent Web Chat Backend (host process)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ndrs
WorkingDirectory=/home/ndrs/hermes-agent-web-chat/backend

Environment=HERMES_HOME=/home/ndrs/.hermes
Environment=HERMES_SRC=/home/ndrs/.hermes/hermes-agent
# Set AUTH_PASSWORD via EnvironmentFile to keep it out of the repo:
#   Environment=AUTH_PASSWORD=your-secret-password
Environment=PORT=11300

# Load sensitive environment variables from a file excluded from git:
# EnvironmentFile=/etc/hermes-agent-web-chat.env

ExecStart=/home/ndrs/.hermes/hermes-agent/venv/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 11300

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
